# script-tools/full/backup_restore/systemd/checkmk-backup-job01.service
[Unit]
Description=CheckMK Job01 Weekly Backup Management (Full 362MB)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/checkmk-tools/script-tools/full/backup_restore/checkmk_manage_job01_weekly.sh
User=root
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
