# script-tools/auto-git-sync.service
[Unit]
Description=Auto Git Sync Service
After=network.target

[Service]
Type=simple
User=marzio
WorkingDirectory=/home/marzio
ExecStart=/bin/bash /home/marzio/checkmk-tools/script-tools/auto-git-sync.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
