# script-tools/full/systemd/checkmk-backup-job00.timer
[Unit]
Description=CheckMK Job00 Daily Backup Timer
Requires=checkmk-backup-job00.service

[Timer]
OnCalendar=daily
OnCalendar=03:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
