# script-tools/full/systemd/checkmk-backup-job00.service
[Unit]
Description=CheckMK Job00 Daily Backup Management (Compressed 1.2MB)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/checkmk-tools/script-tools/full/checkmk_manage_job00_daily.sh
User=root
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
