# install/checkmk-installer/templates/systemd/ydea-ticket-monitor.service
[Unit]
Description=Ydea Ticket Monitor
After=network.target

[Service]
Type=oneshot
User=monitoring
ExecStart=/opt/ydea-toolkit/ydea-ticket-monitor.sh
StandardOutput=journal
StandardError=journal
