# install/checkmk-installer/templates/systemd/ydea-ticket-monitor.timer
[Unit]
Description=Ydea Ticket Monitor Timer
Requires=ydea-ticket-monitor.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=1min

[Install]
WantedBy=timers.target
